# deploy/paper.service
[Unit]
Description=Paper one-time secret sharing service
Documentation=https://github.com/rselbach/paper
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=paper
Group=paper
WorkingDirectory=/var/lib/paper
ExecStart=/usr/local/bin/paper
Restart=on-failure
RestartSec=5s

Environment=PAPER_ADDR=127.0.0.1:8080
Environment=PAPER_DB=/var/lib/paper/paper.db
Environment=PAPER_PUBLIC_ORIGIN=https://paper.example.com
Environment=PAPER_SECRET_TTL_HOURS=168
Environment=PAPER_CLEANUP_INTERVAL_MINUTES=60
Environment=PAPER_MAX_SECRET_BYTES=65536

NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/var/lib/paper

[Install]
WantedBy=multi-user.target
